# deploy/tadmor-reseed.timer
[Unit]
Description=Nightly tadmor demo database reseed

[Timer]
OnCalendar=*-*-* 04:17:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
